# .config/systemd/user/restic-backup.service
[Unit]
Description=Backup system

[Service]
ExecStart=systemd-inhibit ~/.local/bin/restic-backup
